# systemd/sss.service
[Unit]
Description=Wiki Sprechstunden Bot

[Service]
Type=simple
ExecStart=/usr/local/bin/wiki-scripts/sss.py
