# systemd/upgrade.timer
[Unit]
Description=Wiki Upgrade Timer

[Timer]
# Daily at 1:07:21
OnCalendar=*-*-* 01:07:21

[Install]
WantedBy=timers.target
